# deploy/oracle/systemd/ai-stock-arena-api.service
[Unit]
Description=AI Stock Arena FastAPI service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/ai-stock-arena/current
EnvironmentFile=/etc/ai-stock-arena/ai-stock-arena.env
ExecStart=/usr/bin/env bash /opt/ai-stock-arena/current/scripts/linux/run-api.sh
Restart=always
RestartSec=10
StandardOutput=append:/var/log/ai-stock-arena/api.log
StandardError=append:/var/log/ai-stock-arena/api-error.log

[Install]
WantedBy=multi-user.target
