# deploy/oracle/systemd/ai-stock-arena-scheduler.service
[Unit]
Description=AI Stock Arena runtime scheduler
After=network.target ai-stock-arena-api.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/ai-stock-arena/current
EnvironmentFile=/etc/ai-stock-arena/ai-stock-arena.env
ExecStart=/usr/bin/env bash /opt/ai-stock-arena/current/scripts/linux/run-scheduler.sh
Restart=always
RestartSec=10
TimeoutStopSec=30
KillMode=mixed
MemoryAccounting=true
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/var/log/ai-stock-arena/scheduler.log
StandardError=append:/var/log/ai-stock-arena/scheduler-error.log

[Install]
WantedBy=multi-user.target
